# src/components/demoApp.service
[Unit]
Description=myApp
After=syslog.target

[Service]
ExecStart=java /opt/demoApp/deno.jar

[Install]
WantedBy=multi-user.target
